# deploy/clawlodge.service
[Unit]
Description=ClawLodge Next.js App
After=network.target

[Service]
Type=simple
WorkingDirectory=/var/www/clawlodge
Environment=NODE_ENV=production
Environment=PORT=3001
Environment=CLAWLODGE_DATA_DIR=/var/lib/clawlodge
EnvironmentFile=/var/www/clawlodge/.env.production
ExecStart=/bin/sh -lc 'exec ./node_modules/.bin/next start -p "${PORT}"'
ExecStop=/bin/kill -s SIGINT $MAINPID
Restart=always
RestartSec=5
TimeoutStopSec=20
KillMode=mixed

[Install]
WantedBy=multi-user.target
