# deploy/clawlodge.service
[Unit]
Description=ClawLodge Next.js App
After=network.target

[Service]
Type=simple
WorkingDirectory=/var/www/clawlodge
Environment=NODE_ENV=production
Environment=PORT=3001
Environment=CLAWLODGE_DATA_DIR=/var/lib/clawlodge
EnvironmentFile=/var/www/clawlodge/.env.production
ExecStart=/usr/bin/npm run start
ExecStartPost=/bin/sh -lc 'sleep 2; /usr/bin/curl -fsS "http://127.0.0.1:3001/api/v1/lobsters?page=1&pageSize=12" >/dev/null || true; /usr/bin/curl -fsS "http://127.0.0.1:3001/" >/dev/null || true'
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
